# mosquitto.copy/usr/rc.local
#!/bin/sh
set +e ## ignore errors

echo $HOME > /root/home

exit 0