# nas/usr/rc.local
#!/bin/sh

# if BRIDGED wait for interface up
[ ${BRIDGED} ] && pipework --wait

# user integrity
nas-user integrity


exit 0
