# nas/usr/rc.local
#!/bin/sh

#if BRIDGED wait for interface up
[ ${BRIDGED} ] && pipework --wait

# user´s integrity
nas-user integrity


exit 0
