# backup/usr/rc.local
#!/bin/sh
set +e ## ignore errors

# custom rc.local
/srv/conf.d/rc.local

# duplicity & duply
mkdir -p /tmp/cache
ln -s /tmp/cache $HOME/.cache
ln -s /srv/conf.d/duply $HOME/.duply


# update tools on backgroud
nohup setup_tool update duplicity &
nohup setup_tool update duply &


exit 0
