# backup/usr/rc.local
#!/bin/sh
set +e ## ignore errors

# custom rc.local
/srv/conf.d/rc.local


# update tools on backgroud
nohup setup_tool update duplicity &
nohup setup_tool update duply &


# nfs_automount
## https://github.com/vwal/nfs_automount
ln -s /srv/conf.d/nfs/nfs_automount.conf /etc/nfs_automount.conf


exit 0
