# mariadb/usr/rc.local
#!/bin/sh

DATADIR=/srv/mariadb/mysql
[ "$(ls -A ${DATADIR} 2>/dev/null)" ] || {
	mkdir -p ${DATADIR}
	chown -R mysql.mysql ${DATADIR}
	# Set datadir
	mysql_install_db --user=mysql --datadir=${DATADIR}
	# Set user/password admin and grant login from everywhere
#	echo "GRANT ALL ON *.* TO admin@'%' IDENTIFIED BY 'admin' WITH GRANT OPTION; FLUSH PRIVILEGES" | mysql
}

exit 0
