# youtrack/usr/rc.local
#!/bin/sh
set +e ## ignore errors

# custom rc.local
/srv/conf.d/rc.local


# update tools
nohup setup_tool update java &
nohup setup_tool update jetty &
nohup setup_tool update tomcat &
nohup setup_tool update youtrack &

exit 0
